# DataBase/Docker/cron/crontab
* * * * * echo "[SUPERCRONIC-TEST] $(date)" >> /var/log/cron.log
* * * * * /scripts/backup.sh
* * * * * /scripts/cleanup_script.sh
