# DataBase/Docker/cron/crontab
* * * * * echo "[SUPERCRONIC-TEST] $(date)" >> /var/log/cron.log
* * * * * /backup.sh
