# DataBase/Docker/cron/crontab
# Запускать каждый месяц, 1-го числа, в 3:00
0 3 1 * * psql -h db -U aragami -d voluntarydb -f /sql/clean_old_posts.sql > /proc/1/fd/1 2>&1
